# roles/python/templates/backend.service.j2
[Unit]
Description=Backend service for {{ domain_name }}
After=network.target

[Service]
User={{ backend_user }}
Group={{ backend_user }}
WorkingDirectory={{ backend_dir }}
Environment="PATH={{ backend_dir }}/venv/bin"
Environment="PYTHONUNBUFFERED=1"
StandardOutput=append:/var/log/{{ domain_name }}-backend.log
StandardError=append:/var/log/{{ domain_name }}-backend.log
ExecStart={{ backend_dir }}/venv/bin/python {{ backend_dir }}/main.py
Restart=always
RestartSec=10
StartLimitInterval=60
StartLimitBurst=3

[Install]
WantedBy=multi-user.target 